# .config/autostart/autorun.desktop
[Desktop Entry]
Type=Application
Name=autorun
Comment=Script to exec scripts that are intended to run on the background
Exec=autorun
